# debian/ntpbox-display-splash.service
[Unit]
Description=NTPBox display splash screen

[Service]
User=root
ExecStart=/usr/bin/ntpbox-display splash

[Install]
WantedBy=multi-user.target